# deploy/systemd/trade-post.service
[Unit]
Description=交易系统盘后同步（pull + post + git push）

[Service]
Type=oneshot
WorkingDirectory=PATH_TO_TRADESYSTEM
ExecStart=/bin/bash scripts/sync_data.sh
