# deploy/launchd/com.alyx.tradesystem.recommend-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.recommend-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/recommend-runner.sh</string>
        <string>weekly</string>
    </array>

    <!-- 周日 20:00（Weekday: 0=周日，也可写 7） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>20</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-recommend-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-recommend-weekly.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
