# deploy/launchd/com.alyx.tradesystem.trend-leader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Sleep policy: 趋势主升观察清单错过可接受（盘后只读非交易决策任务）；macOS 休眠期间不触发，可接受 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.trend-leader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/trend-leader-runner.sh</string>
    </array>

    <!-- 工作日 21:30（排在 volume-watch 21:00 + sector-correlation 21:15 之后，
         确保 daily_volume_concentration 已落库，主线池不回退）。Weekday: 1=周一 ... 5=周五。
         非交易日（节假日）仍会触发，但当日无涨停/无集中度，报告无新入池。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-trend-leader.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-trend-leader.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
