# deploy/launchd/com.alyx.tradesystem.board-break.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Sleep policy: 观察类清单，错过可接受 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.board-break</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/board-break-runner.sh</string>
    </array>

    <!-- 工作日 21:20（排在 volume-watch 21:00 + sector-correlation 21:15 之后、
         trend-leader 21:30 之前——主线板块归属取 daily_volume_concentration 当日快照）。
         Weekday: 1=周一 ... 5=周五。非交易日（节假日）仍会触发，任务内自动跳过。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>20</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>20</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-board-break.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-board-break.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
