# deploy/launchd/com.alyx.tradesystem.research-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Deprecated: 研报速读生产定时已迁移到 Codex 自动化「每日慧博研报速读（Computer Use）」。
     不再安装本 launchd plist；保留文件仅用于历史排障参考。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.research-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/research-digest-runner.sh</string>
    </array>

    <!-- 每天 22:00 触发；runner 仅在 A 股交易日或 A 股交易日前一天继续执行。 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>22</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-research-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-research-digest.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
