# deploy/launchd/com.alyx.tradesystem.macro-flash.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Sleep policy: 宏观快讯速读错过可接受(非交易决策任务);休眠漏跑用 macro-flash run --date 手动补跑 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.macro-flash</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/macro-flash-runner.sh</string>
    </array>

    <!-- 工作日 16:30(盘后当日快讯)+ 周日 22:00(周末汇总,runner 内切 54h 回溯)。
         Weekday: 0=周日, 1=周一 ... 5=周五。节假日照常触发(宏观新闻不休市)。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-macro-flash.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-macro-flash.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
