# deploy/launchd/com.alyx.tradesystem.research-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Sleep policy: 研报速读错过可接受（非交易决策任务）；macOS 休眠期间不触发，可接受，不配 pmset 唤醒。
     06:42 比 today-pre 07:00 / recommend 07:10 早，天然错峰，均为 I/O 短任务无资源争用。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.research-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/research-digest-runner.sh</string>
    </array>

    <!-- 工作日 06:42（盘前，早于 07:00/07:10 错峰）。Weekday: 1=周一 ... 5=周五。
         A股取最近交易日、美股按美东窗口；非交易日触发也会显式标注无数据，不报错。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>42</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>42</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>42</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>42</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>42</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-research-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-research-digest.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
