# deploy/launchd/com.alyx.tradesystem.daily-leaders.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Sleep policy: 每日最票候选确认稿错过可接受（盘后确认辅助任务）；macOS 休眠期间不触发，可接受 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alyx.tradesystem.daily-leaders</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alyx/tradeSystem/deploy/launchd/daily-leaders-runner.sh</string>
    </array>

    <!-- 工作日 22:30（排在 post-market 派生任务之后，生成复盘第 5 步候选确认稿）。
         Weekday: 1=周一 ... 5=周五。非交易日（节假日）仍会触发，任务内按数据可用性处理。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tradesystem-daily-leaders.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tradesystem-daily-leaders.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
